# rootfs/etc/systemd/system/cron.timer
[Unit]
Description=cron

[Timer]
OnBootSec=10m
OnUnitInactiveSec=24h
WakeSystem=true

[Install]
WantedBy=timers.target
