# rootfs/etc/systemd/system/prune.timer
[Unit]
Description=prune

[Timer]
OnBootSec=10m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
